# shadowsocks-kcptun/files/etc/s6/kcptun/run
#!/bin/execlineb -P
multisubstitute
{
  importas -D 444    KCPTUN_PORT  KCPTUN_PORT
  importas -D 443    SERVER_PORT  SERVER_PORT
  importas -D none   KCPTUN_CRYPT KCPTUN_CRYPT
  importas -D 1200   KCPTUN_MTU   KCPTUN_MTU
  importas -D normal KCPTUN_MODE  KCPTUN_MODE
  importas -D 46     KCPTUN_DSCP  KCPTUN_DSCP
}
foreground
{
  s6-svwait /usr/local/etc/s6/ss
}
kcptun-server -l :$KCPTUN_PORT -t 127.0.0.1:$SERVER_PORT --crypt $KCPTUN_CRYPT --mtu $KCPTUN_MTU --nocomp --mode $KCPTUN_MODE --dscp $KCPTUN_DSCP
