# sks/files/etc/s6/sks-log-clean/run
#!/bin/execlineb -P
loopwhilex foreground {
    cd /var/lib/sks/KDB db_archive -d
    cd /var/lib/sks/PTtree db_archive -d
    sleep 604800
}
