# shadowsocks-kcptun/files/etc/s6/ss/run
#!/bin/execlineb -P
multisubstitute
{
  importas -D 0.0.0.0                 SERVER_ADDR SERVER_ADDR
  importas -D 443                     SERVER_PORT SERVER_PORT
  importas -D foo!bar                 PASSWORD    PASSWORD
  importas -D xchacha20-ietf-poly1305 METHOD      METHOD
  importas -D 300                     TIMEOUT     TIMEOUT
  importas -D 8.8.8.8                 DNS_ADDR    DNS_ADDR
}
ss-server -s $SERVER_ADDR -p $SERVER_PORT -k $PASSWORD -m $METHOD -t $TIMEOUT -d $DNS_ADDR --fast-open -u
