# sks/files/etc/s6/sks-recon/run
#!/bin/sh
s6-svwait /etc/s6/sks-db
sks -stdoutlog -basedir /var/lib/sks recon
