# sks/files/etc/s6/sks-db/run
#!/bin/execlineb -P
sks -stdoutlog -basedir /var/lib/sks db
