# sks/files/etc/s6/sks-db/run
#!/bin/sh
exec sks -stdoutlog -basedir /var/lib/sks db
