# shadowsocks-kcptun/files/s6/ss/run
#!/bin/execlineb -P
multisubstitute
{
  import -D 0.0.0.0 SERVER_ADDR
  import -D 443 SERVER_PORT
  import -D foo!bar PASSWORD
  import -D chacha20 METHOD
  import -D 300 TIMEOUT
  import -D 8.8.8.8 DNS_ADDR
}
ss-server -s $SERVER_ADDR -p $SERVER_PORT -k $PASSWORD -m $METHOD -t $TIMEOUT -d $DNS_ADDR --fast-open -u -A
