# sshd/files/s6/huproxy/run
#!/bin/sh
set -ex

port="${HUPROXY_PORT:-8888}"

exec /usr/local/bin/huproxy -listen ":$port"
