# sks-full/files/etc/s6/sks-recon/run
#!/bin/sh
s6-svwait /usr/local/etc/s6/sks-db
exec sks -stdoutlog -basedir /var/lib/sks recon
