# bind9/files/s6/bind9/run
#!/bin/sh
cd /etc/bind || exit 1
[ -e /usr/local/bind9/etc/rndc.key ] || /usr/local/bind9/sbin/rndc-confgen -a
exec /usr/local/bind9/sbin/named -g
