# bind9/files/s6/bind9/run
#! /bin/execlineb -P
cd /etc/bind /usr/local/bind9/sbin/named -g
