# bind9/files/s6/webhook/run
#! /bin/execlineb -P
/usr/local/bin/webhook -hooks /etc/webhook/hooks.json -verbose
