# bind9/files/s6/webhook/run
#!/bin/sh
exec /usr/local/bin/webhook -hooks /etc/webhook/hooks.json -verbose
