# sks-full/files/etc/s6/cron/run
#!/bin/execlineb -P
crond -f -c /usr/local/etc/crontabs
