# sks/files/etc/s6/cron/run
#!/bin/sh
exec crond -f -c /usr/local/etc/crontabs
