# sks-full/files/etc/s6/caddy/run
#!/bin/sh
export CADDYPATH=/var/lib/sks/caddy
exec caddy -conf /var/lib/sks/caddy/Caddyfile -agree
