# sks-full/files/etc/s6/caddy/run
#!/bin/execlineb -P
export CADDYPATH /var/lib/sks/caddy
/usr/bin/caddy -conf /var/lib/sks/caddy/Caddyfile -agree
